# devel/etc/catkin/profile.d/python_orocos_kdl_site_packages.zsh
export PYTHONPATH=/home/songhat/catkin_ws/devel/lib/python3/dist-packages:$PYTHONPATH
